# ansible/physical/roles/btrfs-balance/files/btrfs-balance@.timer
[Unit]
Description=Periodically start a btrfs balance on %f

[Timer]
OnCalendar=Monthly
Persistent=True

[Install]
WantedBy=timers.target
